# packages/google-cloud-cli/google-cloud-cli.container
[Container]
Image=gcr.io/google.com/cloudsdktool/google-cloud-cli:latest
AutoUpdate=registry
Exec=sleep infinity
Mount=type=bind,source=%E/gcloud,target=/root/.config/gcloud

[Quadlet]
DefaultDependencies=false
